# systemd/captain-fin-drive-sync.service
[Unit]
Description=Captain Fin submitted reports Google Drive sync

[Service]
Type=oneshot
WorkingDirectory=/home/alexey/GitHub/captain-fin
ExecStart=/home/alexey/GitHub/captain-fin/scripts/server-drive-sync.py
